# 01-razer-udev.rules.template
# UDEV rules for razer mice
#
# We just send a rescan command to razerd.
# It will pick up new mice and forward information to the applications.

ACTION=="add", SUBSYSTEM=="usb_device", SYSFS{idVendor}=="1532", RUN+="$$INSTDIR$$/bin/razercfg --scan"
ACTION=="remove", SUBSYSTEM=="usb_device", RUN+="$$INSTDIR$$/bin/razercfg --scan"

