# razerd.service.template
[Unit]
Description=Razer device state daemon
Requires=local-fs.target
After=local-fs.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/razerd -f

[Install]
WantedBy=multi-user.target
